# systemd/runtipi-companion-backup-daily.timer
[Unit]
Description=Daily Runtipi Companion backup

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
Unit=runtipi-companion-backup@daily.service

[Install]
WantedBy=timers.target
